#  server-file.service
[Unit]
Description=File Server with Token Authentication1
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/server-file --server 0.0.0.0 --port 8080 --upload-dir /var/lib/server-file/uploads
WorkingDirectory=/var/lib/server-file
Restart=always
RestartSec=10

# 环境变量
Environment=LOG_LEVEL=info

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=server-file

# 安全配置
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target